# systemd/auto-maintain.service
[Unit]
Description=AutoMaintain checker

[Service]
ExecStart=/etc/monitor-scripts/auto-maintain.sh
