# systemd/bt-check.service
[Unit]
Description=Bluetooth status checker

[Service]
ExecStart=/etc/monitor-scripts/bt.sh
