# systemd/bt-service-check.service
[Unit]
Description=BlueTooth Advertise service start/stop flag checker

[Service]
ExecStart=/etc/monitor-scripts/service-bt.sh
