# systemd/clean-memory.timer
[Unit]
Description=Run the clean memory job every hour.

[Timer]
OnBootSec=1h
OnUnitActiveSec=1h
AccuracySec=1m

[Install]
WantedBy=timers.target
