# systemd/clear-blockchain-check.service
[Unit]
Description=Clear Blockchain checker

[Service]
ExecStart=/etc/monitor-scripts/clear-blockchain.sh
