# systemd/clean-memory.service
[Unit]
Description=Memory cleanup 

[Service]
ExecStart=/etc/monitor-scripts/clean-memory.sh