# systemd/clear-blockchain-check.timer
[Unit]
Description=Run the clear blockchain status check every 15 seconds.

[Timer]
OnBootSec=15
OnUnitActiveSec=15
AccuracySec=1ms

[Install]
WantedBy=timers.target
