# systemd/cpu-check.service
[Unit]
Description=CPU Utilization checker

[Service]
ExecStart=/etc/monitor-scripts/cpu.sh
