# systemd/external-ip-check.service
[Unit]
Description=External IP checker

[Service]
ExecStart=/etc/monitor-scripts/external-ip.sh
