# systemd/fastsync-check.service
[Unit]
Description=Check the status of fastsync

[Service]
ExecStart=/etc/monitor-scripts/fastsync.sh
