# systemd/gps-check.service
[Unit]
Description=GPS status checker

[Service]
ExecStart=/etc/monitor-scripts/gps.sh
