# systemd/infoheight-check.service
[Unit]
Description=Check the current info height of the miner

[Service]
ExecStart=/etc/monitor-scripts/info-height.sh
