# systemd/install-dashboard.service
[Unit]
Description=Finishes off the installation of the dashboard

[Service]
ExecStart=/etc/monitor-scripts/install-dashboard.sh
