# systemd/local-ip-check.service
[Unit]
Description=Local IP checker

[Service]
ExecStart=/etc/monitor-scripts/local-ip.sh
