# systemd/miner-service-check.service
[Unit]
Description=Miner service start/stop flag checker

[Service]
ExecStart=/etc/monitor-scripts/service-miner.sh
