# systemd/miner-check.service
[Unit]
Description=Miner status checker

[Service]
ExecStart=/etc/monitor-scripts/miner.sh
