# systemd/miner-version-check.service
[Unit]
Description=Miner version checker

[Service]
ExecStart=/etc/monitor-scripts/miner-version-check.sh
