# systemd/password-check.service
[Unit]
Description=Password reset checker

[Service]
ExecStart=/etc/monitor-scripts/password-set-check.sh
