# systemd/peer-list-check.service
[Unit]
Description=Peer list checker

[Service]
ExecStart=/etc/monitor-scripts/peer-list.sh
