# systemd/pf-check.service
[Unit]
Description=Packet Forwarder checker

[Service]
ExecStart=/etc/monitor-scripts/packet-forwarder.sh
