# systemd/pf-service-check.service
[Unit]
Description=Packet Forwarder service start/stop flag checker

[Service]
ExecStart=/etc/monitor-scripts/service-pf.sh
