# systemd/pubkeys-check.service
[Unit]
Description=Scrapes public key

[Service]
ExecStart=/etc/monitor-scripts/pubkeys.sh
