# systemd/pubkeys-check.timer
[Unit]
Description=Scrape the pubkey every hour

[Timer]
OnBootSec=3600
OnUnitActiveSec=3600
AccuracySec=1ms

[Install]
WantedBy=timers.target
