# systemd/pubkeys-check.timer
[Unit]
Description=Scrape the pubkey every half hour

[Timer]
OnBootSec=1800
OnUnitActiveSec=1800
AccuracySec=1ms

[Install]
WantedBy=timers.target
