# systemd/reboot-check.service
[Unit]
Description=Reboot flag checker

[Service]
ExecStart=/etc/monitor-scripts/reboot-check.sh
