# systemd/sn-check.service
[Unit]
Description=Serial Number checker

[Service]
ExecStart=/etc/monitor-scripts/sn-check.sh
