# systemd/temp-check.service
[Unit]
Description=Temperature checker

[Service]
ExecStart=/etc/monitor-scripts/temp.sh
