# systemd/update-check.service
[Unit]
Description=Update checker

[Service]
ExecStart=/etc/monitor-scripts/update-check.sh
