# systemd/update-miner-check.service
[Unit]
Description=Miner updater checker

[Service]
ExecStart=/etc/monitor-scripts/miner-update.sh
