# systemd/wifi-check.service
[Unit]
Description=WiFi Status checker

[Service]
ExecStart=/etc/monitor-scripts/wifi-check.sh
