# systemd/wifi-service-check.service
[Unit]
Description=WiFi service start/stop flag checker

[Service]
ExecStart=/etc/monitor-scripts/service-wifi.sh
