# systemd/wifi-config-check.service
[Unit]
Description=New WiFi Config checker

[Service]
ExecStart=/etc/monitor-scripts/wifi-config.sh
